# releng/airootfs/etc/udev/rules.d/devify.rules
SUBSYSTEM=="hid", ACTION=="add|remove", \
RUN+="/usr/local/bin/devify \"$env{HID_NAME}\" \"$env{ACTION}\"", GOTO="end_lookup"
